# crontab
@reboot /usr/local/bin/tmux new -s mobile -d
