# crontab
*/5 * * * * $HOME/bin/duck.sh >/dev/null 2>&1
